# xe-mailserver-conf/lib/systemd/system/mailserver.target
[Unit]
Description=Dovecot, Exim4 and all other services for a mailserver

[Install]
WantedBy=multi-user.target
